# cart.service
#!/bin/bash

[Unit]
Description = Cart Service
[Service]
User=roboshop
Environment=REDIS_HOST=redis.nagendrababu.online
Environment=CATALOGUE_HOST=catalogue.nagendrababu.online
Environment=CATALOGUE_PORT=8080
ExecStart=/bin/node /app/server.js
SyslogIdentifier=cart

[Install]
WantedBy=multi-user.target